# recipes-nexus/nexus/files/storage-daemon.service
[Unit]
Description=Nexus Storage Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/storage-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
